# pylons/pylons-monitor.service
[Unit]
Description=Additional Pylons Load Monitor

[Service]
Type=oneshot
ExecStart=%h/.claude/pylons/pylons-monitor.sh
