# minerals/minerals-monitor.service
[Unit]
Description=Not Enough Minerals RAM Monitor

[Service]
Type=oneshot
ExecStart=%h/.claude/minerals/minerals-monitor.sh
